# crontab
*/5 * * * * /usr/local/bin/python /code/LeadFuzion/web_app/order_complition_mail.py
